# doorbell-notifier.service
[Unit]
Description=Doorbell MQTT Listener
After=network.target

[Service]
Type=simple
ExecStart=%h/.local/bin/doorbell-notifier
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
